# deploy/ansible/roles/services/templates/gistbook.conf
#!upstart
#
# Copy me to /etc/init/
# This installs a daemon as a system level call and ensures the process is consistently restarted on error.
# Manual start, stop, and restart respected.
#

description "Daemon for nodejs gistbook server"

start on startup
stop on shutdown
respawn

env NODE_ENV={{node_env}}
script
 /usr/bin/node {{app_base_path}}server/app.js >> /var/log/gistbook 2>&1
end script
